# jobs/crontab
MAILTO=""
0 22 * * Mon . /home/dotcloud/env/bin/activate && /home/dotcloud/current/jobs/knowtify.py | mail -s 'Knowtify API Weekly Update' alewisohn@sitterfied.com,jnanscawen@sitterfied.com
