# jobs/crontab
00 18 * * Mon ~/current/jobs/knowtify | mail -s 'Knowtify API Update Run' alewisohn@sitterfied.com,jnanscawen@sitterfied.com
