# ansible/roles/confd/files/confd.upstart
description "confd daemon"

start on started consul
stop on stopped consul

respawn
respawn limit 10 5

console log

pre-start script
  /bin/sleep 8
end script

exec /usr/local/bin/confd
